# templates/default/services/cpan-mirror.conf.erb
#!upstart
description "Update CPAN mirror"

start on started cpan-mirror-shim
stop on stopping cpan-mirror-shim

task
respawn

setuid <%= node.cpan_mirror.user %>

env USER=<%= node.cpan_mirror.user %>
env HOME=/home/<%= node.cpan_mirror.user %>

exec rsync -av --delete cpan-rsync.perl.org::CPAN <%= node.cpan_mirror.data_dir %>
